# etc/xxvd
#!/bin/sh
#
# xxvd          Start/Stop the xxvd daemon.
#
# chkconfig: 2345 90 11
# description: xxvd is a telnet/http config       \
# server, for ther DVB Settop box vdr program.       \
#                                                    \
#                                                    \
#                                                    .
# processname: xxvd
# config: xxvd.cfg
#
# Need to setup :
# Define FOLDER directory, ( $FOLDER/bin/xxvd )
FOLDER="/opt/XXV"

# How verbose should log file (0 ... 5)
VERBOSE="3"

# Run as noprivileged user, else set empty
RUNAS="vdr"

# Which translation is used
LANGUAGE="de_DE@euro";

# Adjust some moduls directories
MODPATH="$FOLDER/lib/XXV/MODULES"
CONFIG="$FOLDER/etc/xxvd.cfg"
HTMLDIR="$FOLDER/skins/"
DOCUDIR="$FOLDER/doc/"
PODDIR="$FOLDER/doc/"
CONTRIB="$FOLDER/contrib"
NEWSMODS="$FOLDER/lib/XXV/OUTPUT/NEWS"
NEWSTMPL="$FOLDER/share/news"

PIDFILE="/var/run/xxvd.pid"
LOGFILE="/var/log/xxvd.log"

# Create options from 
OPTIONS="--configfile $CONFIG \
         --logfile $LOGFILE \
         --pidfile $PIDFILE \
         --docudir $DOCUDIR \
         --poddir $PODDIR \
         --htmldir $HTMLDIR \
         --contrib $CONTRIB \
         --newsmods $NEWSMODS \
         --newstmpl $NEWSTMPL \
         --moduledir $MODPATH \
         --verbose $VERBOSE"

# See how we were called.
PROG="bin/xxvd"

start() {
  echo -n "Start $PROG: "

  if [ -r "$PIDFILE" ]; then
    PID=$(cat "$PIDFILE")
    if [ -r "/proc/$PID/cmdline" ]; then
      grep "xxvd" /proc/$PID/cmdline >/dev/null 2>&1 \
        || rm -f "$PIDFILE"
    else
      rm -f "$PIDFILE"
    fi
  fi

  su - $RUNAS -c "export LANG="$LANGUAGE";cd $FOLDER;nice -n 10 ./$PROG $OPTIONS"
}

stop() {
  echo -n "Stop $PROG: "
  su - $RUNAS -c "export LANG="$LANGUAGE";cd $FOLDER;./$PROG --kill $OPTIONS"
}

restart() {
  stop
  start
}

case "$1" in
  start)
    start
    exit 0;
    ;;
  stop)
    stop
    exit 0;
    ;;
  restart)
    restart
    exit 0;
    ;;
  *)
    echo $"Usage: $0 {start|stop|status|restart}"
    exit 1
esac

exit $?
